# agent.service
[Unit]
Description=Autonomous Site Agent
After=network.target

[Service]
Type=simple
User=nosyt
WorkingDirectory=/home/nosyt/.openclaw/workspace
ExecStart=/home/nosyt/.openclaw/workspace/autonomous-agent.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
